# debian/cron.d
* * * * * root ganeti-failover-gateway-node
