# debian/cron.d
* * * * * root ganeti-failover-gateway-node
* * * * * root sleep 10 && ganeti-failover-gateway-node
* * * * * root sleep 20 && ganeti-failover-gateway-node
* * * * * root sleep 30 && ganeti-failover-gateway-node
* * * * * root sleep 40 && ganeti-failover-gateway-node
* * * * * root sleep 50 && ganeti-failover-gateway-node
